# etc/init.d/egpu
#!/sbin/openrc-run

command="/usr/local/bin/egpu-switcher switch auto"
