# etc/local.d/nyan.start
#!/bin/bash
hdparm -B 127 /dev/sdb
fsck -a /dev/mapper/weedbanger
mount /mnt/weedbanger/

